# ports/openssh/.footprint
drwxr-sr-x	root/users	usr/   
drwxr-sr-x	root/users	usr/lib/   
drwxr-sr-x	root/users	usr/lib/openssh/   
-rwxr-xr-x	root/users	usr/lib/openssh/sftp-server   
-rwxr-xr-x	root/users	usr/lib/openssh/ssh-pkcs11-helper   
-rws--x--x	root/users	usr/lib/openssh/ssh-keysign   
drwxr-sr-x	root/users	usr/share/   
drwxr-sr-x	root/users	usr/share/man/   
drwxr-sr-x	root/users	usr/share/man/man1/   
lrwxrwxrwx	root/users	usr/share/man/man1/slogin.1.gz -> ./ssh.1.gz 
-rw-r--r--	root/users	usr/share/man/man1/scp.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/ssh.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/ssh-add.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/ssh-agent.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/sftp.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/ssh-keyscan.1.gz   
-rw-r--r--	root/users	usr/share/man/man1/ssh-keygen.1.gz   
drwxr-sr-x	root/users	usr/share/man/man5/   
-rw-r--r--	root/users	usr/share/man/man5/moduli.5.gz   
-rw-r--r--	root/users	usr/share/man/man5/sshd_config.5.gz   
-rw-r--r--	root/users	usr/share/man/man5/ssh_config.5.gz   
drwxr-sr-x	root/users	usr/share/man/man8/   
-rw-r--r--	root/users	usr/share/man/man8/ssh-keysign.8.gz   
-rw-r--r--	root/users	usr/share/man/man8/sshd.8.gz   
-rw-r--r--	root/users	usr/share/man/man8/ssh-pkcs11-helper.8.gz   
-rw-r--r--	root/users	usr/share/man/man8/sftp-server.8.gz   
drwxr-sr-x	root/users	usr/sbin/   
-rwxr-xr-x	root/users	usr/sbin/sshd   
drwxr-sr-x	root/users	usr/bin/   
lrwxrwxrwx	root/users	usr/bin/slogin -> ./ssh 
-rwxr-xr-x	root/users	usr/bin/sftp   
-rwxr-xr-x	root/users	usr/bin/ssh-keygen   
-rwxr-xr-x	root/users	usr/bin/ssh-keyscan   
-rwxr-xr-x	root/users	usr/bin/ssh-agent   
-rwxr-xr-x	root/users	usr/bin/ssh-add   
-rwxr-xr-x	root/users	usr/bin/scp   
-rwxr-xr-x	root/users	usr/bin/ssh   
drwxr-sr-x	root/users	etc/   
drwxr-sr-x	root/users	etc/rc.d/   
-rwxr-xr-x	root/users	etc/rc.d/sshd   
drwxr-sr-x	root/users	etc/ssh/   
-rw-r--r--	root/users	etc/ssh/sshd_config   
-rw-r--r--	root/users	etc/ssh/ssh_config   
-rw-r--r--	root/users	etc/ssh/moduli   
drwxr-sr-x	root/users	var/   
drwxr-sr-x	root/users	var/lib/   
drwxr-sr-x	root/sys	var/lib/sshd/   
