# unsupported/packages/sysklogd/.footprint
drwxr-xr-x	root/root	var/   
drwxr-xr-x	root/root	var/log/   
-rw-r-----	root/root	var/log/debug   
-rw-r-----	root/root	var/log/messages   
-rw-r-----	root/root	var/log/mail   
-rw-r-----	root/root	var/log/cron   
-rw-r-----	root/root	var/log/auth   
-rw-r-----	root/root	var/log/kernel   
drwxr-xr-x	root/root	usr/   
drwxr-xr-x	root/root	usr/sbin/   
-rwxr-xr-x	root/root	usr/sbin/rotatelog   
-rwx------	root/root	usr/sbin/klogd   
-rwx------	root/root	usr/sbin/syslogd   
drwxr-xr-x	root/root	usr/share/   
drwxr-xr-x	root/root	usr/share/man/   
drwxr-xr-x	root/root	usr/share/man/man8/   
-rw-r--r--	root/root	usr/share/man/man8/klogd.8.gz   
-rw-r--r--	root/root	usr/share/man/man8/syslogd.8.gz   
-rw-r--r--	root/root	usr/share/man/man8/sysklogd.8.gz   
drwxr-xr-x	root/root	usr/share/man/man5/   
-rw-r--r--	root/root	usr/share/man/man5/syslog.conf.5.gz   
drwxr-xr-x	root/root	etc/   
-rw-r--r--	root/root	etc/syslog.conf   
drwxr-xr-x	root/root	etc/rc.d/   
-rwxr-xr-x	root/root	etc/rc.d/sysklogd   
drwxr-xr-x	root/root	etc/cron/   
drwxr-xr-x	root/root	etc/cron/weekly/   
-rwxr-xr-x	root/root	etc/cron/weekly/syslog   
