# unsupported/packages/dcron/crontab
#
# /etc/crontab - crond(8) configuration
#
# The format of this file:
# Min Hour Day Month DayOfWeek Command

00 20 * * *	/usr/sbin/runjobs /etc/cron/daily
05 20 * * mon	/usr/sbin/runjobs /etc/cron/weekly
10 20 1 * *	/usr/sbin/runjobs /etc/cron/monthly

# End of file
